# userspace_program/driver_communication_service/systemd/printer_lamp_driver_service.service
[Unit]
Description=Printer lamp driver interaction service

[Service]
ExecStart=/usr/lib/printer_lamp/driver_interaction
Restart=on-failure
StartLimitBurst=0

[Install]
WantedBy=dbus.service
